# root/etc/services.d/cronicle/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/bin/node /opt/cronicle/lib/main.js start
