# root/etc/services.d/cronicle/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /opt/cronicle/bin/control.sh start
